# pi-agent/systemd/css-auto-update.service
[Unit]
Description=CSS Signage Auto-Update Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/css
ExecStart=/bin/bash -c 'git -C /opt/css -c safe.directory=/opt/css fetch origin main && git -C /opt/css -c safe.directory=/opt/css reset --hard origin/main && bash /opt/css/pi-agent/scripts/setup-log-rotation.sh && systemctl restart css-agent'
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
