# pi-agent/systemd/css-daily-reboot.service
[Unit]
Description=CSS Signage Daily Reboot

[Service]
Type=oneshot
ExecStart=/sbin/shutdown -r now
